# etc/supervisord.d/dormer-receiver.ini
[group:dormer-receiver]
programs=2003,2103

[program:2003]
command=/opt/dormer-receiver/venv/bin/python3 app.py 2003
directory=/opt/dormer-receiver/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/dormer-receiver-stderr.log
stdout_logfile=/var/log/supervisor/dormer-receiver-stdout.log

[program:2103]
command=/opt/dormer-receiver/venv/bin/python3 app.py 2103
directory=/opt/dormer-receiver/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/dormer-receiver-stderr.log
stdout_logfile=/var/log/supervisor/dormer-receiver-stdout.log
